# txtdot.service
[Unit]
Description=txtdot proxy server
After=network.target

[Service]
User=txtdot
Group=txtdot
WorkingDirectory=/home/txtdot/src/dist
ExecStart=/usr/bin/node ./src/app.js
Restart=always

[Install]
WantedBy=default.target
